# home/.config/systemd/user/xkb.service
[Unite]
Description=Assign custom keymap

[Service]
ExecStart=bash %h/.config/xkb/enable-keymap.sh
